# supervisord.d/consul.conf
[program:consul]
command=/usr/local/bin/consul agent -config-dir=/data/consul/config
user=root
