# supervisord.d/consul-template.conf
[program:consul-template]
command=/usr/local/bin/consul-template -config=/data/consul-template/config.hcl
user=root
